# scripts/com.docmind.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.docmind.web</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/kashika/dev/RAG/scripts/run-web.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/kashika/dev/RAG/apps/web</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>3</integer>
  <key>StandardOutPath</key>
  <string>/tmp/rag-web.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/rag-web.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>NODE_ENV</key>
    <string>development</string>
  </dict>
</dict>
</plist>
